# scripts/polymarket-bot.service
[Unit]
Description=NachoMarket Polymarket Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/nachomarket
Environment=PYTHONPATH=/home/ubuntu/nachomarket
EnvironmentFile=-/home/ubuntu/nachomarket/.env

# Pre-flight: verificar que la IP no está geobloqueada
ExecStartPre=/home/ubuntu/nachomarket/venv/bin/python scripts/check_geo.py

ExecStart=/home/ubuntu/nachomarket/venv/bin/python src/main.py
Restart=on-failure
RestartSec=30

# Safety: restart on failure but not too aggressively
StartLimitInterval=300
StartLimitBurst=5

# Resource limits (Free Tier: 1GB RAM)
MemoryMax=800M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nachomarket

[Install]
WantedBy=multi-user.target
